# files/packetads.service
[Unit]
Description=PacketLoss Advertisements Bot Service
After=network.target

[Service]
User=packetloss
WorkingDirectory=/home/packetloss/packetads
ExecStart=/home/packetloss/.local/share/pnpm/pnpm run start
Restart=always
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=packetads

[Install]
WantedBy=multi-user.target